# config/includes.chroot/etc/xdg/autostart/seconionis-mate.desktop
[Desktop Entry]
Type=Application
Name=Seconionis
Exec=seconionis-mate
OnlyShowIn=MATE;
StartupNotify=false
Terminal=false
