# config/includes.chroot/etc/init.d/jason
#!/bin/bash

### BEGIN INIT INFO
# Provides:          jason
# Required-Start:
# Required-Stop:
# Should-Start:
# Default-Start:  2 3 4 5
# Default-Stop:  0 1 6
# Short-Description: Transparent Proxy through TOR.
### END INIT INFO

# Developed by root0emir

export BLUE='\033[1;94m'
export GREEN='\033[1;92m'
export RED='\033[1;91m'
export RESETCOLOR='\033[1;00m'

# Destinations you don't want routed through Tor
TOR_EXCLUDE="192.168.0.0/16 172.16.0.0/12 10.0.0.0/8"

# The UID Tor runs as
# change it if, starting tor, the command 'ps -e | grep tor' returns a different UID
TOR_UID="debian-tor"

# Tor's TransPort and SocksPort
TOR_TRANS_PORT="9040"
TOR_SOCKS_PORT="9050"

# DNS Port
DNS_PORT="5353"

# Is resolvconf installed?
if [[ ! -z `dpkg -l |grep resolvconf` ]]; then 
	resolvconf_support=true;
else
	resolvconf_support=false;
fi


function init {
	echo -e -n " $GREEN*$BLUE killing dangerous applications$RESETCOLOR\n"
	killall -q chrome dropbox iceweasel skype icedove thunderbird firefox chromium xchat transmission deluge pidgin pidgin.orig
	
	echo -e -n " $GREEN*$BLUE cleaning some dangerous cache elements"
	bleachbit -c adobe_reader.cache chromium.cache chromium.current_session chromium.history elinks.history emesene.cache epiphany.cache firefox.url_history flash.cache flash.cookies google_chrome.cache google_chrome.history  links2.history opera.cache opera.search_history opera.url_history &> /dev/null
}

function disable_ipv6() {
    echo -e "\n$GREEN[$BLUE i$GREEN ]$BLUE Stopping IPv6 services:$RESETCOLOR\n"

    # Remove current IPv6 settings if it exists
    if [ -f /etc/sysctl.d/98-jason-privacy.conf ]; then
        rm /etc/sysctl.d/98-jason-privacy.conf
    fi

    # Add comprehensive IPv6 disable configuration
    cat << EOF >> /etc/sysctl.d/98-jason-privacy.conf
# Disable IPv6 - Added by jason
net.ipv6.conf.all.disable_ipv6 = 1 #jason
net.ipv6.conf.default.disable_ipv6 = 1 #jason
net.ipv6.conf.lo.disable_ipv6 = 1 #jason
EOF

    # Apply settings
    if ! sysctl -p /etc/sysctl.d/98-jason-privacy.conf > /dev/null; then
        echo -e "$RED Error applying sysctl settings$RESETCOLOR"
        return 1
    fi

    # Verify IPv6 is disabled
    if [ "$(cat /proc/sys/net/ipv6/conf/all/disable_ipv6)" != "1" ]; then
        echo -e "$RED Failed to disable IPv6$RESETCOLOR"
        return 1
    fi

    echo -e "$GREEN IPv6 successfully disabled$RESETCOLOR"
    return 0
}

function enable_ipv6() {
    echo -e "\n$GREEN[$BLUE i$GREEN ]$BLUE Restarting Network services:$RESETCOLOR\n"

    # Remove the config file if it exists
    if [ -f /etc/sysctl.d/98-jason-privacy.conf ]; then
        rm -f /etc/sysctl.d/98-jason-privacy.conf
		echo -e "$GREEN Removed jason sysctl config$RESETCOLOR"
    fi

    # Reload all sysctl settings
    if ! sysctl --system > /dev/null; then
        echo -e "$RED Error reloading sysctl settings$RESETCOLOR"
        return 1
    fi

    # Restart network services
    service network-manager force-reload > /dev/null 2>&1
    service nscd start > /dev/null 2>&1
    service dnsmasq start > /dev/null 2>&1

}


function ip {
	echo -e "\nMy ip is:\n"
	sleep 1
	curl "https://api.ipify.org" # Decent provider serving 30b requests/month, and doesn't block TOR
	echo -e "\n\n----------------------------------------------------------------------"
}

function start {
	# Make sure only root can run this script
	if [ $(id -u) -ne 0 ]; then
		echo -e -e "\n$GREEN[$RED!$GREEN] $RED This script must be run as root$RESETCOLOR\n" >&2
		exit 1
	fi
	
	# Check defaults for Tor
	grep -q -x 'RUN_DAEMON="yes"' /etc/default/tor
	if [ $? -ne 0 ]; then
		echo -e "\n$GREEN[$RED!$GREEN]$RED Please add the following to your /etc/default/tor and restart service:$RESETCOLOR\n" >&2
		echo -e "$BLUE#----------------------------------------------------------------------#$RESETCOLOR"
		echo -e 'RUN_DAEMON="yes"'
		echo -e "$BLUE#----------------------------------------------------------------------#$RESETCOLOR\n"
		exit 1
	fi	
	
	disable_ipv6

	echo -e "\n$GREEN[$BLUE i$GREEN ]$BLUE Starting anonymous mode:$RESETCOLOR\n"
	
	if [ ! -e /var/run/tor/tor.pid ]; then
		echo -e " $RED*$BLUE Tor is not running! $GREEN starting it $BLUE for you$RESETCOLOR\n" >&2
		service network-manager force-reload > /dev/null 2>&1
		killall dnsmasq > /dev/null 2>&1
		killall nscd > /dev/null 2>&1	
		service tor start
		sleep 3
	fi
	
	# Check if Tor config has required settings
	if ! grep -q "^TransPort" /etc/tor/torrc; then
		echo "TransPort 9040" >> /etc/tor/torrc
	fi
	if ! grep -q "^DNSPort" /etc/tor/torrc; then
		echo "DNSPort 5353" >> /etc/tor/torrc
	fi
	if ! grep -q "^ControlPort" /etc/tor/torrc; then
		echo "ControlPort 9051" >> /etc/tor/torrc
	fi
	
	if ! [ -f /etc/network/iptables.rules ]; then
		iptables-save > /etc/network/iptables.rules
		echo -e " $GREEN*$BLUE Saved iptables rules$RESETCOLOR\n"
	fi
	
	iptables -F
	iptables -t nat -F
	
	if [ "$resolvconf_support" = false ] 
	then
		cp /etc/resolv.conf /etc/resolv.conf.bak
		echo -e 'nameserver 127.0.0.1\nnameserver 1.1.1.1' > /etc/resolv.conf
		echo -e " $GREEN*$BLUE Modified resolv.conf to use Tor DNS"
	else
		cp /etc/resolvconf/resolv.conf.d/head{,.bak}
        echo -e 'nameserver 127.0.0.1\nnameserver 1.1.1.1' >> /etc/resolvconf/resolv.conf.d/head
        echo -e " $GREEN*$BLUE Modified resolvconf to use Tor DNS$RESETCOLOR\n"
        resolvconf -u
	fi

	# set iptables nat
	iptables -t nat -A OUTPUT -m owner --uid-owner $TOR_UID -j RETURN
	iptables -t nat -A OUTPUT -p udp --dport 53 -j REDIRECT --to-ports $DNS_PORT
	iptables -t nat -A OUTPUT -p tcp --dport 53 -j REDIRECT --to-ports $DNS_PORT
	
	#resolve .onion domains mapping 10.192.0.0/10 address space
	iptables -t nat -A OUTPUT -p tcp -d 10.192.0.0/10 -j REDIRECT --to-ports $TOR_TRANS_PORT
	iptables -t nat -A OUTPUT -p udp -d 10.192.0.0/10 -j REDIRECT --to-ports $TOR_TRANS_PORT
	
	#exclude local addresses
	for NET in $TOR_EXCLUDE 127.0.0.0/9 127.128.0.0/10; do
		iptables -t nat -A OUTPUT -d $NET -j RETURN
	done
	
	#redirect all other output through TOR
	iptables -t nat -A OUTPUT -p tcp --syn -j REDIRECT --to-ports $TOR_TRANS_PORT
	iptables -t nat -A OUTPUT -p udp -j REDIRECT --to-ports $TOR_TRANS_PORT
	iptables -t nat -A OUTPUT -p icmp -j REDIRECT --to-ports $TOR_TRANS_PORT
	
	#accept already established connections
	iptables -A OUTPUT -m state --state ESTABLISHED,RELATED -j ACCEPT
	
	#exclude local addresses
	for NET in $TOR_EXCLUDE 127.0.0.0/8; do
		iptables -A OUTPUT -d $NET -j ACCEPT
	done
	
	#allow only tor output
	iptables -A OUTPUT -m owner --uid-owner $TOR_UID -j ACCEPT
	iptables -A OUTPUT -j REJECT

	echo -e "$GREEN *$BLUE All traffic was redirected through Tor$RESETCOLOR\n"
	echo -e "$GREEN[$BLUE i$GREEN ]$BLUE You are under Jason tunnel$RESETCOLOR\n"
	
	# Show current IP
	echo -e "$GREEN *$BLUE Your current Tor IP is:$RESETCOLOR"
	curl -s "https://api.ipify.org"
	echo -e "\n"
}

function stop {
	# Make sure only root can run our script
	if [ $(id -u) -ne 0 ]; then
		echo -e "\n$GREEN[$RED!$GREEN] $RED This script must be run as root$RESETCOLOR\n" >&2
		exit 1
	fi
	echo -e "\n$GREEN[$BLUE i$GREEN ]$BLUE Stopping anonymous mode:$RESETCOLOR\n"

	iptables -F
	iptables -t nat -F
	echo -e " $GREEN*$BLUE Deleted all iptables rules\n$RESETCOLOR"
	
	if [ -f /etc/network/iptables.rules ]; then
		iptables-restore < /etc/network/iptables.rules
		rm /etc/network/iptables.rules
		echo -e " $GREEN*$BLUE Iptables rules restored"
	fi
	
	# restore DNS settings
	if [ "$resolvconf_support" = false ] 
	then
		if [ -e /etc/resolv.conf.bak ]; then
			rm /etc/resolv.conf
			cp /etc/resolv.conf.bak /etc/resolv.conf
		fi
	else
		mv /etc/resolvconf/resolv.conf.d/head{.bak,}
		resolvconf -u
	fi
	
	service tor stop
	
	echo -e "\n$GREEN[$BLUE i$GREEN ]$BLUE Reenabling IPv6 services:$RESETCOLOR\n"

	# reenable IPv6 services
	enable_ipv6
	
	echo -e " $GREEN*$BLUE Anonymous mode stopped$RESETCOLOR\n"
}

function change {
    echo -e "\n$GREEN[$BLUE i$GREEN ]$BLUE Changing Tor identity:$RESETCOLOR\n"
    
    # Restart Tor service
    service tor restart
    
    # Wait for Tor to establish new circuits
    echo -e "$GREEN *$BLUE Waiting for new Tor circuits...$RESETCOLOR"
    sleep 5
    
    # Show new IP
    echo -e "$GREEN *$BLUE Your new IP address is:$RESETCOLOR"
    curl -s "https://api.ipify.org"
    echo -e "\n"
}

function status {
	service tor status
}

case "$1" in
	start)
		init
		start
	;;
	stop)
		init
		stop
	;;
	change)
		change
	;;
	status)
		status
	;;
	myip|ip)
		ip
	;;
	restart)
		$0 stop
		sleep 1
		$0 start
	;;
   *)
echo -e "
Jason Privacy Tool
	Usage:
	$RED┌──[$GREEN$USER$YELLOW@$BLUE`hostname`$RED]─[$GREEN$PWD$RED]
	$RED└──╼ \$$GREEN"" jason $RED{$GREEN""start$RED|$GREEN""stop$RED|$GREEN""restart$RED|$GREEN""change$RED""$RED|$GREEN""status$RED""}
	
	$RED start$BLUE -$GREEN Start system-wide anonymous
		  tunneling under TOR proxy through iptables	  
	$RED stop$BLUE -$GREEN Reset original iptables settings
		  and return to clear navigation
	$RED restart$BLUE -$GREEN Combines \"stop\" and \"start\" options
	$RED change$BLUE -$GREEN Changes identity restarting TOR
	$RED status$BLUE -$GREEN Check if Jason is working properly
	$RED myip$BLUE -$GREEN Show your current IP address $RED(RISKY)$RESETCOLOR
	
$RESETCOLOR" >&2
exit 1
;;
esac

echo -e $RESETCOLOR
exit 0
